# py-utils/src/utils/message_bus/cortx_message_bus.service
[Unit]
Description=cortx rest server
After=network.target

[Service]
User=root
ExecStart=/opt/seagate/cortx/utils/bin/utils_server
Restart=always

[Install]
WantedBy=multi-user.target